# com.andrewfraser.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.andrewfraser.dashboard</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/montymac/.openclaw/workspace/content-pipeline/venv/bin/python3</string>
        <string>/Users/montymac/.openclaw/workspace/content-pipeline/app.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/montymac/.openclaw/workspace/content-pipeline</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/dashboard.out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/dashboard.err.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
